# .config/systemd/user/offlineimap-quick.target
[Unit]
Description=OfflineIMAP one-time quick-sync for all enabled accounts
