# .config/systemd/user/offlineimap.target
[Unit]
Description=All OfflineIMAP sync accounts
Requires=sockets.target
