# crontab
# vim:set tw=0:o
LC_ALL=C.UTF-8
# m h  dom mon dow   command
0    0 * * * [ -e ~/.logrotate.conf ] && /usr/sbin/logrotate -s ~/.logrotate.state ~/.logrotate.conf
